# scripts/com.knuask.codex-enrichment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.knuask.codex-enrichment</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/noseongjun/Library/Application Support/KNU-Ask/codex_ingestion_worker.py</string>
        <string>--watch</string>
        <string>--interval</string>
        <string>1</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/noseongjun/Library/Application Support/KNU-Ask</string>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- macOS 개인정보 보호 경로 밖의 권한 600 파일에서 토큰을 읽는다. -->
        <key>KNUASK_ENV_FILE</key>
        <string>/Users/noseongjun/Library/Application Support/KNU-Ask/worker.env</string>
        <key>KNUASK_API_BASE</key>
        <string>http://localhost:8080/api</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noseongjun/Library/Logs/KNU-Ask-codex-worker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noseongjun/Library/Logs/KNU-Ask-codex-worker-error.log</string>
</dict>
</plist>
